# playbooks/roles/collectd_server/templates/collectd_proxy.service
[Unit]
Description=collectd_proxy
After=network.target

[Service]
ExecStart=/usr/bin/python /opt/collectd_proxy/server.py
User=collectd_proxy
Group=collectd_proxy
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
